# build/ports/transmission/files/transmission.in
#!/bin/sh
# Copyright (c) 2006-2008 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: transmission
# REQUIRE: DAEMON
# KEYWORD: shutdown
# XQUERY: -i "count(//bittorrent/enable) > 0" -o "0" -b
# RCVAR: transmission

. /etc/rc.subr
. /etc/util.subr
. /etc/configxml.subr

name="transmission"
rcvar=`set_rcvar`

load_rc_config "${name}"

# Custom commands
start_precmd="start_precmd"

# Defaults
transmission_enable=${transmission_enable:-"NO"}
transmission_remote_user=${transmission_remote_user:-"admin"}
transmission_user=${transmission_user:-"transmission"}
transmission_suopt=${transmission_suopt:-"-l"}

command="/usr/local/bin/${name}-daemon"
command_args="--no-blocklist --auth --username ${transmission_remote_user}"

start_precmd()
{
	local _command_args

	# Create command args line.
	_command_args=`/usr/local/bin/xml sel -t -m "//bittorrent" \
		-v "concat('--password ',password,' --download-dir &quot;',downloaddir,'&quot; --port ',port)" \
		-i "string-length(configdir) > 0" -v "concat(' --config-dir &quot;',configdir,'&quot;')" -b \
		-o " --allowed *" \
		${configxml_file} | /usr/local/bin/xml unesc`

	command_args="${command_args} ${_command_args}"
}

run_rc_command "$1"
