# build/ports/transmission/files/transmission.in
#!/bin/sh
# Copyright (c) 2006-2009 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: transmission
# REQUIRE: DAEMON
# KEYWORD: shutdown
# XQUERY: -i "count(//bittorrent/enable) > 0" -o "0" -b
# RCVAR: transmission

. /etc/rc.subr
. /etc/util.subr
. /etc/configxml.subr

name="transmission"
rcvar=`set_rcvar`

load_rc_config "${name}"

# Custom commands
start_precmd="start_precmd"

# Defaults
transmission_enable=${transmission_enable:-"NO"}
transmission_remoteuser=${transmission_remoteuser:-"admin"}
transmission_allowed=${transmission_allowed:-"*.*.*.*"}
transmission_noblocklist=${transmission_noblocklist:-"YES"}  # To check if we really need noblocklist 
transmission_blocklist=${transmission_blocklist:-"NO"}
transmission_user=${transmission_user:-"transmission"}
transmission_suopt=${transmission_suopt:-"-l"}

command="/usr/local/bin/${name}-daemon"
command_args="--allowed ${transmission_allowed}"

start_precmd()
{
	local _command_args

	# Create command args line.
	_command_args=`/usr/local/bin/xml sel -t -m "//bittorrent" \
		-i "count(authrequired) > 0" \
			-o " --auth" \
			-o " --username ${transmission_remoteuser}" \
			-v "concat(' --password ',password)" \
		-b \
		-i "count(authrequired) = 0" -o " --no-auth" -b \
		-v "concat(' --download-dir &quot;',downloaddir,'&quot; --port ',port)" \
		-i "string-length(configdir) > 0" -v "concat(' --config-dir &quot;',configdir,'&quot;')" -b \
		${configxml_file} | /usr/local/bin/xml unesc`

	if checkyesno transmission_noblocklist; then
		_command_args="${_command_args} --no-blocklist"
	fi
	if checkyesno transmission_blocklist; then
		_command_args="-b ${_command_args}"
	fi
	command_args="${command_args} ${_command_args}"
}

run_rc_command "$1"
