# etc/rc.d/inadyn
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: inadyn
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr
. /etc/rc.subrconf

name="inadyn"
rcvar=`set_rcvar`
command="/usr/local/bin/${name}"
command_args="--input_file /var/etc/${name}.conf"

load_rc_config $name
create_conf_file $name
run_rc_command "$1"
