# etc/rc.d/iscsi_initiator
#!/bin/sh
# Copyright (c) 2007-2008 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: iscsi_initiator
# REQUIRE: NETWORKING
# KEYWORD: shutdown
# XQUERY: -i "count(//iscsiinit/vdisk) > 0" -o "0" -b
# RCVAR: iscsi_initiator

. /etc/rc.subr
. /etc/configxml.subr

name="iscsi_initiator"
rcvar=`set_rcvar`

# Custom commands
start_precmd="start_precmd"
start_cmd="start_cmd"
stop_cmd="stop_cmd"

# Defaults
iscsi_initiator_enable="${iscsi_initiator_enable:=NO}"
iscsi_initiator_conf=${iscsi_initiator_conf:-"/var/etc/iscsi.conf"}
iscsi_initiator_port=${iscsi_initiator_port:-"3260"}
command="/sbin/iscontrol"
command_args="-c ${iscsi_initiator_conf}"

start_precmd()
{
	# Load kernel module
	load_kld -e iscsi_initiator iscsi_initiator || return 1

	# Disable debugging
	/sbin/sysctl debug.iscsi_initiator=0 > /dev/null

	# Create config file
	/usr/local/bin/xml sel -t \
		-o "port = ${iscsi_initiator_port}" -n \
		-m "//iscsiinit/vdisk" \
			-v "concat(name,' {')" -n \
			-v "concat('  targetaddress = ',targetaddress)" -n \
			-v "concat('  targetname = ',targetname)" -n \
			-v "concat('  initiatorname = ',initiatorname)" -n \
			-o "}" -n \
		-b \
		${configxml_file} | /usr/local/bin/xml unesc > ${iscsi_initiator_conf}
}

start_cmd()
{
	local _name

	echo "Starting ${name}."

	/usr/local/bin/xml sel -t -m "//iscsiinit/vdisk" \
		-v "name" \
		-i "position() != last()" -n -b \
		${configxml_file} | /usr/local/bin/xml unesc | \
		while read _name; do
			eval ${command} ${command_args} -n ${_name} > /dev/null
		done
}

stop_cmd()
{
	echo "Stopping ${name}."
	/usr/bin/killall iscontrol
}

load_rc_config ${name}
run_rc_command "$1"
