# etc/rc.d/iscsi_target
#!/bin/sh

# PROVIDE: iscsi_target
# REQUIRE: NETWORKING
# XQUERY: -i "count(//iscsitarget/vdisk) > 0" -o "0" -b
# RCVAR: iscsi_target

. /etc/rc.subr
. /etc/configxml.subr

name="iscsi_target"
rcvar=`set_rcvar`

# Custom commands
start_precmd="create_conf"

# Defaults
iscsi_target_enable=${iscsi_target_enable:-"NO"}
iscsi_target_config="/var/etc/iscsi/targets"
command="/usr/local/bin/iscsi-target"
command_args="-f ${iscsi_target_config}"
pidfile="/var/run/iscsi-target.pid"

create_conf()
{
	# Clear existing file.
	cat /dev/null > ${iscsi_target_config}

	# Iterate over all configured disks.
	/usr/local/bin/xml sel -t -m "//iscsitarget/vdisk" \
		-v "concat('extent',position()-1,' /mnt/',sharename,'/iscsi_target',position()-1,' 0 ',size,'MB')" -n \
		-v "concat('target',position()-1,' rw extent',position()-1,' ',network)" -n \
		${configxml_file} >> ${iscsi_target_config}
}

load_rc_config ${name}
run_rc_command "$1"
