# etc/rc.d/mdnsresponder
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: mdnsresponder
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr
. /etc/conf.subr

name="mdnsresponder"
rcvar=`set_rcvar`
mdnsresponder_config=${mdnsresponder_config:-"/var/etc/${name}.conf"}
required_files=${mdnsresponder_config}
command="/usr/local/sbin/mDNSResponderPosix"
command_args="-b -f ${mdnsresponder_config}"
pidfile="/var/run/mDNSResponder.pid"
start_precmd="start_precmd"

start_precmd()
{
	# Create .conf file
	create_conf_file ${name}
}

load_rc_config ${name}
run_rc_command "$1"
