# etc/rc.d/msntp
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: msntp
# REQUIRE: DAEMON

. /etc/rc.subr
. /etc/configxml.subr

name="msntp"
rcvar=`set_rcvar`

# Custom commands
start_cmd="start_cmd"
stop_cmd="stop_cmd"

# Defaults
msntp_enable="${msntp_enable:=NO}"
pidfile="/var/run/${name}.pid"
scriptpidfile="/var/run/runmsntp.pid"

start_cmd()
{
	local _updateinterval _timeservers

	echo "Starting ${name}."

	_updateinterval=`/usr/local/bin/xml sel -t -v "//system/time-update-interval" ${configxml_file}`
	_timeservers=`/usr/local/bin/xml sel -t -v "//system/timeservers" ${configxml_file}`

	nohup /usr/local/bin/runmsntp.sh "${scriptpidfile}" "${pidfile}" "${_updateinterval}" "${_timeservers}" > /dev/null 2>&1 &
}

stop_cmd()
{
	echo "Stopping ${name}."

	# Kill runmsntp script.
	/bin/kill -s TERM `/bin/cat ${scriptpidfile}`
	rm -r ${scriptpidfile}

	# Kill msntp.
	/bin/kill -s TERM `/bin/cat ${pidfile}`
	rm -r ${pidfile}
}

load_rc_config ${name}
run_rc_command "$1"
