# etc/rc.d/pureftpd
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: pureftpd
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

pureftpd_enable=${pureftpd_enable-"NO"}

. /etc/rc.subr
. /etc/conf.subr

name="pureftpd"
rcvar=`set_rcvar`

load_rc_config ${name}

command="/usr/local/sbin/pure-ftpd"
pureftpd_config=${pureftpd_config:-"/var/etc/pure-ftpd.conf"}
required_files=${pureftpd_config}
pidfile="/var/run/pure-ftpd.pid"
start_precmd="start_precmd"

start_precmd()
{
	# Create .conf file
	create_conf_file ${name}

	# Convert config file into command args (pure-ftpd does not support config files).
	for _command_arg in `cat ${pureftpd_config}`; do
		command_args="${command_args} ${_command_arg}"
	done
}

run_rc_command "$1"
