# etc/rc.d/rcconf
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: rcconf

. /etc/rc.subr
. /etc/conf.subr

name="rcconf"
stop_cmd=":"
start_cmd=":"
start_precmd="start_precmd"

start_precmd()
{
	# Create .conf file
	create_conf_file ${name}
}

load_rc_config ${name}
run_rc_command "$1"
