# etc/rc.d/ushare
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: ushare
# REQUIRE: DAEMON
# KEYWORD: shutdown
# XQUERY: -i "count(//upnp/enable) > 0" -o "0" -b
# RCVAR: ushare

. /etc/rc.subr
. /etc/configxml.subr

name="ushare"
rcvar=`set_rcvar`

# Custom commands
start_precmd="start_precmd"
stop_postcmd="stop_postcmd"

# Defaults
ushare_enable=${ushare_enable:-"NO"}
command="/usr/local/bin/${name}"
command_args="--daemon --override-iconv-err"

start_precmd()
{
	local _if _command_args

	# Get interface used for uPnP.
	_if=`configxml_get "//upnp/if"`

	# Setup a multicast route for UPnP messages.
	/sbin/route add -net 239.0.0.0 -netmask 255.0.0.0 -interface ${_if}

	# Create command args line.
	_command_args=`/usr/local/bin/xml sel -t -m "//upnp" \
		-i "profile[. = 'xbox']" -o " --xbox" -b \
		-i "profile[. = 'dlna']" -o " --dlna" -b \
		-i "count(web) = 0" -o " --no-web" -b \
		-v "concat(' --name=&quot;',name,'&quot;')" \
		-v "concat(' --interface=',if)" \
		-i "string-length(port) > 0" -v "concat(' --port=',port)" -b \
		-m "//upnp/content" -v "concat(' --content=&quot;',.,'&quot;')" -b \
		${configxml_file}`

	command_args="${command_args} ${_command_args}"
}

stop_postcmd()
{
	# Delete multicast route.
	/sbin/route delete -net 239.0.0.0 > /dev/null 2>&1
}

load_rc_config ${name}
run_rc_command "$1"
