# etc/rc.d/ushare
#!/bin/sh
# Copyright © 2007 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: ushare
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr
. /etc/rc.subrconf

name="ushare"
rcvar=`set_rcvar`
ushare_config=${ushare_config:-"/var/etc/${name}.conf"}
required_files=${ushare_config}
command="/usr/local/bin/${name}"
command_args="--daemon --xbox --cfg=${ushare_config}"
stop_postcmd="ushare_stop_postcmd"

ushare_stop_postcmd()
{
	# Delete multicast route.
	/sbin/route delete -net 239.0.0.0 > /dev/null
}

load_rc_config ${name}
create_conf_file ${name}
run_rc_command "$1"
