# build/ports/fuppes/files/fuppes.in
#!/bin/sh
# Copyright (c) 2006-2009 Volker Theile (votdev@gmx.de)
# All rights reserved.

# PROVIDE: fuppes
# REQUIRE: DAEMON
# KEYWORD: shutdown
# XQUERY: -i "count(//upnp/enable) > 0" -o "0" -b
# RCVAR: fuppes

. /etc/rc.subr
. /etc/util.subr
. /etc/configxml.subr

name="fuppes"
rcvar=`set_rcvar`

load_rc_config "${name}"

# Custom commands
start_precmd="fuppes_mkconf"
stop_postcmd="fuppes_poststop"
mkconf_cmd="fuppes_mkconf"
rebuilddb_cmd="fuppes_dboptions"
updatedb_cmd="fuppes_dboptions"
rebuildvcont_cmd="fuppes_dboptions"
extra_commands="mkconf rebuilddb updatedb rebuildvcont"

# Defaults
fuppes_enable=${fuppes_enable:-"NO"}
fuppes_configdir=${fuppes_configdir-"/var/etc"}
fuppes_config=${fuppes_config-"${fuppes_configdir}/${name}.cfg"}
fuppes_vfolder_config=${fuppes_vfolder_config-"/usr/local/share/${name}/vfolder.cfg"}
fuppes_loglevel=${fuppes_loglevel-"1"}
fuppes_logfile=${fuppes_logfile-"/var/log/${name}.log"}
fuppes_port=${fuppes_port-"49152"}
fuppes_local_charset=${fuppes_local_charset-"UTF-8"}
fuppes_use_imagemagick=${fuppes_use_imagemagick-"false"}
fuppes_use_taglib=${fuppes_use_taglib-"true"}
fuppes_use_libavformat=${fuppes_use_libavformat-"true"}

command="/usr/local/bin/${name}d"
command_args="--config-dir ${fuppes_configdir} --config-file ${fuppes_config} --vfolder-config-file ${fuppes_vfolder_config} --log-level ${fuppes_loglevel} --log-file ${fuppes_logfile}"

# Set kill signal to KILL, otherwise fuppes takes too long to terminate.
sig_stop=${sig_stop-"KILL"}

fuppes_mkconf()
{
	local _if _port _name _profile _deviceip _command_args

	# Get interface and port used for uPnP.
	_if=`configxml_get "//upnp/if"`
	_port=`configxml_get "//upnp/port"`
	_name=`configxml_get "//upnp/name"`
	_profile=`configxml_get "//upnp/profile"`
	_deviceip=`configxml_get "//upnp/deviceip"`

	# Set customized port.
	[ -n "${_port}" ] && fuppes_port="${_port}"

	# Set default friendly name if necessary.
	[ -z "${_name}" ] && _name=`get_product_name` # Set default friendly name.

	# Setup a multicast route for UPnP messages.
	/sbin/route add -net 239.0.0.0 -netmask 255.0.0.0 -interface ${_if} >/dev/null 2>&1

	# Create config file.
	cat <<EOF > ${fuppes_config}
<?xml version="1.0" encoding="UTF-8"?>
<fuppes_config version="0.7.2.3">
  <shared_objects>
EOF

	/usr/local/bin/xml sel -t -m "//upnp/content" \
		-v "concat('&lt;dir&gt;',.,'&lt;/dir&gt;')" -n \
		${configxml_file} | /usr/local/bin/xml unesc >> ${fuppes_config}

	cat <<EOF >> ${fuppes_config}
  </shared_objects>
  <network>
    <interface>${_if}</interface>
    <http_port>${fuppes_port}</http_port>
    <allowed_ips>
    </allowed_ips>
  </network>
  <content_directory>
  	<!--a list of possible charsets can be found under: http://www.gnu.org/software/libiconv/-->
    <local_charset>${fuppes_local_charset}</local_charset>
    <!--libs used for metadata extraction when building the database. [true|false]-->
    <use_imagemagick>${fuppes_use_imagemagick}</use_imagemagick>
    <use_taglib>${fuppes_use_taglib}</use_taglib>
    <use_libavformat>${fuppes_use_libavformat}</use_libavformat>
  </content_directory>
  <transcoding>
  	<!--[lame|twolame]-->
    <audio_encoder>lame</audio_encoder>
    <!--[true|false]-->
    <transcode_vorbis>false</transcode_vorbis>
    <transcode_musepack>false</transcode_musepack>
    <transcode_flac>false</transcode_flac>
  </transcoding>
  <device_settings>
<!--"default" settings are inhertied by specific devices and can be overwritten-->
<!--do NOT remove the "default" device settings-->
<!--all new file types have to be added to the default settings-->
<!--adding new file types just to a specific device will have no affect-->
    <device name="default" virtual="default" enabled="false">
<!--specify the maximum length for file names (0 or empty = unlimited)-->
      <max_file_name_length>0</max_file_name_length>
<!--[file|container]-->
      <playlist_style>file</playlist_style>
      <show_childcount_in_title>false</show_childcount_in_title>
      <enable_dlna>false</enable_dlna>
      <transcoding_release_delay>4</transcoding_release_delay>
      <description_values>
        <friendly_name>${_name}</friendly_name>
      </description_values>
      <file_settings>
<!--audio files-->
        <file ext="mp3">
          <type>AUDIO_ITEM</type>
          <mime_type>audio/mpeg</mime_type>
          <dlna>MP3</dlna>
        </file>
        <file ext="m4a">
          <type>AUDIO_ITEM_MUSIC_TRACK</type>
          <mime_type>audio/mp4</mime_type>
          <dlna>MP4A</dlna>
        </file>
        <file ext="ogg">
          <type>AUDIO_ITEM</type>
          <mime_type>application/octet-stream</mime_type>
          <transcode enabled="false">
            <ext>mp3</ext>
            <mime_type>audio/mpeg</mime_type>
            <dlna>MP3</dlna>
            <http_encoding>chunked</http_encoding>
            <decoder>vorbis</decoder>
            <encoder>lame</encoder>
            <bitrate>192</bitrate>
            <samplerate>44100</samplerate>
          </transcode>
        </file>
        <file ext="mpc">
          <type>AUDIO_ITEM</type>
          <mime_type>application/octet-stream</mime_type>
          <transcode enabled="false">
            <ext>mp3</ext>
            <mime_type>audio/mpeg</mime_type>
            <dlna>MP3</dlna>
            <http_encoding>chunked</http_encoding>
            <decoder>musepack</decoder>
            <encoder>lame</encoder>
            <bitrate>192</bitrate>
            <samplerate>44100</samplerate>
          </transcode>
        </file>
        <file ext="wav">
          <type>AUDIO_ITEM</type>
          <mime_type>audio/x-wav</mime_type>
        </file>
        <file ext="flac">
          <type>AUDIO_ITEM</type>
          <mime_type>audio/x-flac</mime_type>
          <transcode enabled="false">
            <ext>mp3</ext>
            <mime_type>audio/mpeg</mime_type>
            <dlna>MP3</dlna>
            <http_encoding>chunked</http_encoding>
            <decoder>flac</decoder>
            <encoder>lame</encoder>
            <bitrate>192</bitrate>
            <samplerate>44100</samplerate>
          </transcode>
        </file>
        <file ext="wma">
          <type>AUDIO_ITEM</type>
          <mime_type>audio/x-ms-wma</mime_type>
          <dlna>WMAFULL</dlna>
        </file>
<!--image files-->
        <file ext="jpg">
          <ext>jpeg</ext>
          <type>IMAGE_ITEM</type>
          <mime_type>image/jpeg</mime_type>
          <convert enabled="false">
<!--<dcraw enabled="true">-q 0</dcraw>-->
            <ext>png</ext>
            <mime_type>image/png</mime_type>
            <height>0</height>
            <width>0</width>
<!--set "greater" to "true" if you only want to resize images greater than "height" or "width"-->
            <greater>false</greater>
<!--set "less" to "true" if you only want to resize images less than "height" or "width"-->
            <less>false</less>
<!--set "less" and "greater" to "false" if you always want to resize-->
          </convert>
        </file>
        <file ext="bmp">
          <type>IMAGE_ITEM</type>
          <mime_type>image/bmp</mime_type>
        </file>
        <file ext="png">
          <type>IMAGE_ITEM</type>
          <mime_type>image/png</mime_type>
        </file>
        <file ext="gif">
          <type>IMAGE_ITEM</type>
          <mime_type>image/gif</mime_type>
        </file>
<!--video files-->
        <file ext="mpg">
          <ext>mpeg</ext>
          <type>VIDEO_ITEM</type>
          <mime_type>video/mpeg</mime_type>
        </file>
        <file ext="mp4">
          <type>VIDEO_ITEM</type>
          <mime_type>video/mp4</mime_type>
        </file>
        <file ext="avi">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-msvideo</mime_type>
        </file>
        <file ext="divx">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-divx</mime_type>
        </file>
        <file ext="wmv">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-ms-wmv</mime_type>
        </file>
        <file ext="vob">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-ms-vob</mime_type>
        </file>
        <file ext="vdr">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-extension-vdr</mime_type>
          <transcode enabled="false">
            <ext>vob</ext>
            <mime_type>video/x-ms-vob</mime_type>
          </transcode>
        </file>
        <file ext="flv">
          <type>VIDEO_ITEM</type>
          <mime_type>application/x-flash-video</mime_type>
        </file>
        <file ext="asf">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-ms-asf</mime_type>
        </file>
        <file ext="mkv">
          <type>VIDEO_ITEM</type>
          <mime_type>video/x-matroska</mime_type>
          <transcode enabled="false">
            <transcoder>ffmpeg</transcoder>
            <ext>mpg</ext>
            <mime_type>video/mpeg</mime_type>
            <video_codec>mpeg2video</video_codec>
            <video_bitrate>1800000</video_bitrate>
            <audio_codec>mp2</audio_codec>
            <audio_samplerate>44100</audio_samplerate>
            <audio_bitrate>192000</audio_bitrate>
          </transcode>
        </file>
        <file ext="iso">
          <type>VIDEO_ITEM</type>
          <mime_type>application/x-iso-image</mime_type>
        </file>
        <file ext="m2ts">
          <type>VIDEO_ITEM</type>
          <mime_type>video/mp4</mime_type>
        </file>
<!--playlists-->
        <file ext="pls">
          <type>PLAYLIST</type>
          <mime_type>audio/x-scpls</mime_type>
        </file>
        <file ext="m3u">
          <type>PLAYLIST</type>
          <mime_type>audio/x-mpegurl</mime_type>
        </file>
      </file_settings>
    </device>
  	<device name="PS3" enabled="false">
      <user_agent>UPnP/1.0 DLNADOC/1.50</user_agent>
      <user_agent>PLAYSTATION 3</user_agent>
<!--<ip></ip>-->
      <enable_dlna>true</enable_dlna>
      <transcoding_release_delay>50</transcoding_release_delay>
      <file_settings>
        <file ext="ogg">
          <type>AUDIO_ITEM_MUSIC_TRACK</type>
          <transcode enabled="true">
            <http_encoding>stream</http_encoding>
          </transcode>
        </file>
        <file ext="avi">
          <type>VIDEO_ITEM</type>
          <mime_type>video/divx</mime_type>
        </file>
        <file ext="vob">
          <type>VIDEO_ITEM</type>
          <mime_type>video/mpeg</mime_type>
        </file>
        <file ext="ts">
          <ext>mpeg</ext>
          <type>VIDEO_ITEM</type>
          <mime_type>video/mpeg</mime_type>
        </file>
        <file ext="mpg">
          <type>VIDEO_ITEM</type>
          <transcode enabled="false">
            <transcoder>ffmpeg</transcoder>
            <mime_type>video/mpeg</mime_type>
            <ext>mpg</ext>
            <video_codec>copy</video_codec>
            <audio_codec>mp2</audio_codec>
            <audio_samplerate>44100</audio_samplerate>
            <audio_bitrate>192000</audio_bitrate>
          </transcode>
        </file>
      </file_settings>
    </device>
    <device name="Microsoft_XBox360" virtual="Microsoft_XBox360" enabled="false">
      <user_agent>Xbox/2.0.\d+.\d+ UPnP/1.0 Xbox/2.0.\d+.\d+</user_agent>
      <user_agent>Xenon</user_agent>
      <xbox360>true</xbox360>
      <show_empty_resolution>true</show_empty_resolution>
      <description_values>
        <friendly_name>${_name} : 1</friendly_name>
        <model_name>Windows Media Connect compatible (%s)</model_name>
        <model_number>2.0</model_number>
      </description_values>
      <file_settings>
        <file ext="mp3"><type>AUDIO_ITEM_MUSIC_TRACK</type></file>
        <file ext="jpg"><type>IMAGE_ITEM_PHOTO</type></file>
        <file ext="avi">
          <type>VIDEO_ITEM</type>
          <mime_type>video/avi</mime_type>
        </file>
        <file ext="asf">
          <mime_type>video/x-ms-wmv</mime_type>
        </file>
      </file_settings>
    </device>
    <device name="Terratec_Noxon_iRadio" virtual="default" enabled="false">
<!--Please enter the address of your Noxon. Automatic detection is impossible because the Noxon does not send a "user-agent" in it's requests-->
<!--<ip></ip>-->
      <playlist_style>container</playlist_style>
      <show_childcount_in_title>true</show_childcount_in_title>
    </device>
    <device name="Telegent_TG100" virtual="default" enabled="false">
      <user_agent>dma/1.0 \(http://www.cybertan.com.tw/\)</user_agent>
      <user_agent>UPnP/1.0 DLNADOC/1.00</user_agent>
      <playlist_style>file</playlist_style>
      <max_file_name_length>101</max_file_name_length>
    </device>
    <device name="DLNA" virtual="default" enabled="false">
      <enable_dlna>true</enable_dlna>
    </device>
    <device name="ZyXEL_DMA1000" virtual="default" enabled="false">
    </device>
    <device name="Helios_X3000" virtual="default" enabled="false">
    </device>
    <device name="DLink_DSM320" virtual="default" enabled="false">
      <description_values>
        <user_agent>redsonic</user_agent>
        <manufacturer_url>redsonic.com</manufacturer_url>
        <model_number>105</model_number>
      </description_values>
    </device>
    <device name="Yamaha_RXN600" virtual="Yamaha_RXN600" enabled="false">
      <user_agent>YAMAHA Network Receiver</user_agent>
      <user_agent>YAMAHA Network Receiver DLNADOC/1.00</user_agent>
      <xbox360>false</xbox360>
      <enable_xms_media_receiver_registrar>true</enable_xms_media_receiver_registrar>
      <enable_url_base>false</enable_url_base>
      <description_values>
        <manufacturer>Microsoft</manufacturer>
        <manufacturer_url>http://www.microsoft.com/</manufacturer_url>
        <model_name>Windows Media Connect</model_name>
        <model_number>2.0</model_number>
        <model_url>http://www.microsoft.com/</model_url>
        <model_description enabled="false" />
        <upc enabled="false" />
        <serial_number enabled="false" />
      </description_values>
    </device>
    <device name="Loewe_Connect" virtual="Loewe_Connect" enabled="false">
      <max_file_name_length>0</max_file_name_length>
      <playlist_style>container</playlist_style>
      <show_childcount_in_title>true</show_childcount_in_title>
      <enable_dlna>true</enable_dlna>
      <transcoding_release_delay>4</transcoding_release_delay>
      <xbox360>false</xbox360>
      <show_empty_resolution>false</show_empty_resolution>
      <user_agent>POSIX DLNADOC/1.00 INTEL_NMPR/2.1 UPnP/1.0 Intel MicroStack/1.0.1868</user_agent>
      <description_values>
        <manufacturer>Loewe Opta GmbH</manufacturer>
        <manufacturer_url>http://www.loewe.de</manufacturer_url>
        <model_name>NMP-1000</model_name>
        <model_number>33393</model_number>
        <model_url>http://www.loewe.de/connect</model_url>
        <model_description enabled="true">Loewe Network Mediaplayer</model_description>
        <upc enabled="false"/>
        <serial_number enabled="false"/>
      </description_values>
    </device>
  </device_settings>
</fuppes_config>
EOF

	# Enable configured device and vfolder profile.
	/usr/local/bin/xml edit -P -u "//fuppes_config/device_settings/device[@name='${_profile}']/@enabled" \
		-v true ${fuppes_config} > ${fuppes_config}.tmp
	/bin/mv ${fuppes_config}.tmp ${fuppes_config}

	/usr/local/bin/xml edit -P -u "//fuppes_vfolder_config/vfolder_layout[@device='${_profile}']/@enabled" \
		-v true ${fuppes_vfolder_config} > ${fuppes_vfolder_config}.tmp
	/bin/mv ${fuppes_vfolder_config}.tmp ${fuppes_vfolder_config}

	# Set device ip.
	if [ -n ${_deviceip} ]; then
		/usr/local/bin/xml edit -s "//fuppes_config/device_settings/device[@name='${_profile}']" \
			-t elem -n ip -v "${_deviceip}" ${fuppes_config} > ${fuppes_config}.tmp
		/bin/mv ${fuppes_config}.tmp ${fuppes_config}
	fi

	# Enable/disable web interface.
	if ! configxml_isset //upnp/web; then
		# Set allowed IP to loopback to disable all external IP's.
		/usr/local/bin/xml edit -s "//fuppes_config/network/allowed_ips" \
			-t elem -n ip -v "127.0.0.1" ${fuppes_config} > ${fuppes_config}.tmp
		/bin/mv ${fuppes_config}.tmp ${fuppes_config}
	fi

	# Enable/disable transcoding
	if configxml_isset //upnp/transcoding; then
		/usr/local/bin/xml edit -P \
			-u "//fuppes_config/transcoding/transcode_vorbis" -v true \
			-u "//fuppes_config/transcoding/transcode_flac" -v true \
			-u "//fuppes_config/device_settings/device/file_settings/file/transcode/@enabled" -v true \
			${fuppes_config} > ${fuppes_config}.tmp
		/bin/mv ${fuppes_config}.tmp ${fuppes_config}
	fi

	# Create command args line.
	_command_args=`/usr/local/bin/xml sel -t -m "//upnp" \
		-v "concat('--database-file &quot;',home,'/${name}.db&quot;')" \
		-i "count(transcoding) > 0" \
			-v "concat(' --temp-dir &quot;',tempdir,'&quot;')" \
		-b \
		${configxml_file} | /usr/local/bin/xml unesc`

	command_args="${command_args} ${_command_args}"
}

fuppes_poststop()
{
	# Delete multicast route.
	/sbin/route delete -net 239.0.0.0 >/dev/null 2>&1
}

# Command to
# - rebuild/update database
# - rebuild virtual container
fuppes_dboptions()
{
	local _if _port _ip _param

	# Get interface and port used for uPnP.
	_if=`configxml_get "//upnp/if"`
	_port=`configxml_get "//upnp/port"`
	_ip=`get_ipaddr inet ${_if}`

	case ${rc_arg} in
		rebuilddb)
			_param="db=rebuild";
			;;
		updatedb)
			_param="db=update";
			;;
		rebuildvcont)
			_param="vcont=rebuild";
			;;
	esac

	/usr/bin/fetch -q -o - http://${_ip}:${_port}/presentation/options.html?${_param} 1>/dev/null 2>&1
}

run_rc_command "$1"
