# build/ports/nut/files/notifyme.sh.in
#!/bin/sh
# Copyright (c) 2008 Volker Theile (votdev@gmx.de)
# All rights reserved.

. /etc/rc.subr
. /etc/configxml.subr

# Defaults
msmtp_config=${msmtp_config:-"/var/etc/msmtp.conf"}
msmtp_msgfile=${msmtp_msgfile:-"/tmp/message"}

# Get recipients
_recipients=`configxml_get "//system/email/from"`

# Create message
/usr/local/bin/xml sel -t \
	-v "concat('From: ',//system/email/from)" -n \
	-v "concat('To: ',//ups/to)" -n \
	-v "concat('Subject: ',//ups/subject)" -n \
	-o "." -n \
	${configxml_file} | /usr/local/bin/xml unesc > ${msmtp_msgfile}

# Save the email message (STDIN) to a file:
echo "$*" >> ${msmtp_msgfile}

# Now email the message to the user at address ADD:
/usr/local/bin/msmtp --file=${msmtp_config} ${_recipients} < ${msmtp_msgfile} 1>/dev/null 2>&1

# Cleanup
/bin/rm ${msmtp_msgfile} 1>/dev/null 2>&1
