# build/ports/nut/files/upsmon-notifycmd.in
#!/bin/sh
# Copyright (c) 2008 Volker Theile (votdev@gmx.de)
# All rights reserved.

. /etc/rc.subr
. /etc/configxml.subr

# Defaults
msmtp_config=${msmtp_config:-"/var/etc/msmtp.conf"}
msmtp_msgfile=${msmtp_msgfile:-"/tmp/message"}

# Get recipients
_recipients=`configxml_get "//ups/email/to" | /usr/bin/tr ";" " "`

# Create message
echo ${_recipients} | awk '{for ( i = NF ; i > 0 ; --i ) printf("To: %s\n",$i)}' > ${msmtp_msgfile}

/usr/local/bin/xml sel -t \
	-v "concat('From: ',//system/email/from)" -n \
	-v "concat('Subject: ',//ups/email/subject)" -n \
	-o "." -n \
	${configxml_file} | /usr/local/bin/xml unesc >> ${msmtp_msgfile}

echo "$*" >> ${msmtp_msgfile}

# Now email the message to the user
/usr/local/bin/msmtp --file=${msmtp_config} ${_recipients} < ${msmtp_msgfile} 1>/dev/null 2>&1

# Cleanup
/bin/rm ${msmtp_msgfile} 1>/dev/null 2>&1
