# build/ports/rconf/files/rconf
#!/bin/sh
#
# Copyright (c) 2007 Volker Theile <votdev@gmx.de>.
# All rights reserved.
#
#set -x

type=${1}
action=${2}

. /etc/rc.subr

usage() {
cat  << FOO
usage: ${0} service enable servicename
       ${0} service disable servicename
       ${0} service remove servicename
       ${0} service rcvar servicename
       ${0} attribute set name value
       ${0} attribute remove name

examples:
       # ${0} service enable sshd
       Set sshd_enable="YES" in /etc/rc.conf.

       # ${0} service disable sshd
       Set sshd_enable="NO" in /etc/rc.conf.

       # ${0} service remove sshd
       Remove sshd_enable from /etc/rc.conf.

       # ${0} service rcvar sshd
       Display what rc.conf variable is used (if any).
       Return 0 if YES, otherwise 1

       # ${0} attribute set hostname xyz
       Set attribute hostname="xyz" in /etc/rc.conf.

       # ${0} attribute remove hostname
       Remove attribute hostname from /etc/rc.conf.
FOO
	exit 1
}

# remove string
# Remove given string from file.
remove() {
	/usr/bin/sed "/${1}*/d" /etc/rc.conf > /etc/rc.conf.tmp
	mv -f /etc/rc.conf.tmp /etc/rc.conf
}

if [ $# -lt 2 ]; then
	usage
	exit 1
fi

case ${type} in
	service)
		shift 2
		services=$*
		for service in $services; do
			case ${action} in
				enable)
					remove "${service}_enable"
					echo "${service}_enable=\"YES\"" >> /etc/rc.conf
					;;
			
				disable)
					remove "${service}_enable"
					echo "${service}_enable=\"NO\"" >> /etc/rc.conf
					;;
			
				remove)
					remove "${service}_enable"
					;;
			
				rcvar)
					load_rc_config ${service}
					if checkyesno "${service}_enable"; then
						exit 0
					else
						exit 1
					fi
					;;
			esac
		done
		;;

	attribute)
		attribute=${3}
		case ${action} in
			set)
				remove "${attribute}"
				value=${4}
				echo "${attribute}=\"${value}\"" >> /etc/rc.conf
				;;

			remove)
				remove "${attribute}"
				;;
		esac
		;;

	*)
		usage
		exit 1
		;;
esac

exit 0
