# build/ports/rconf/files/rconf
#!/bin/sh
#
# Copyright (c) 2007 Volker Theile <votdev@gmx.de>.
# All rights reserved.
#
#set -x

cmd=${1}
service=${2}

. /etc/rc.subr

usage() {
cat  << FOO
usage: ${0} enable servicename
       ${0} disable servicename
       ${0} remove servicename
       ${0} rcvar servicename

examples:
       # ${0} enable sshd
       Set sshd_enable="YES" in /etc/rc.conf.

       # ${0} disable sshd
       Set sshd_enable="NO" in /etc/rc.conf.

       # ${0} remove sshd
       Delete sshd_enable in /etc/rc.conf.
       
       # ${0} rcvar sshd
       Display what rc.conf variable is used (if any).
       Return 0 if YES, otherwise 1
FOO
	exit 1
}

remove() {
	for _line in `cat /etc/rc.conf`; do
		if [ ${service} != ${_line%_enable=*} ]; then
			echo ${_line} >> /etc/rc.conf.tmp
		fi
	done
	mv -f /etc/rc.conf.tmp /etc/rc.conf
}

if [ $# -lt 2 ]; then
	usage
	exit 1
fi

case ${cmd} in
	enable)
		remove
		echo "${service}_enable=\"YES\"" >> /etc/rc.conf
		;;

	disable)
		remove
		echo "${service}_enable=\"NO\"" >> /etc/rc.conf
		;;

	remove)
		remove
		;;

	rcvar)
		load_rc_config ${service}
		if checkyesno "${service}_enable"; then
			exit 0
		else
			exit 1
		fi
		;;

	*)
		usage
		exit 1
		;;
esac

exit 0
